# payara/payara.conf
[program:payara]
command=/bin/bash -c "/opt/payara41/bin/asadmin start-domain"

